# rag/dockerize.bak2/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /rag/dockerize/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx